# files/supervisord/conf.d/foreman-proxy.conf
[program:foreman-proxy]
user=foreman
command=ruby /usr/share/foreman-proxy/bin/smart-proxy