# files/supervisord/conf.d/puppet.conf
[program:puppetmaster]
user=puppet
command=puppet master --no-daemonize